# squashfs-root/etc/uci-defaults/luci-i18n-spanish
#!/bin/sh
uci batch <<-EOF
	set luci.languages.es=Espanol
        commit luci
EOF

