# squashfs-root/etc/hotplug.d/iface/81-speedtest
#!/bin/sh
[ "$ACTION" == "ifup" ] && [ "$INTERFACE" == "wan" ] && {
    /usr/bin/oversea_download_xml download_xml
}
